# config/systemd/user/hydroxide.service
[Unit]
Description=Hydroxide

[Service]
ExecStart=/home/heywoodlh/.local/bin/hydroxide serve

[Install]
WantedBy=multi-user.target
